# build/copyfiles/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=service ssh start
#autostart=true
#utorestart=true

[program:wg]
command=wg-quick up wg0
#autostart=true
#utorestart=true

[program:squid]
command=service squid start
#autostart=true
#utorestart=true
